# meta-lenovo/meta-hr855xg2/recipes-phosphor/lenovo-sensors/files/service_files/xyz.openbmc_project.fpgasensor.service
[Unit]
Description=FPGA Sensor
StopWhenUnneeded=false
Requires=logrotate.service xyz.openbmc_project.Logging.IPMI.service fpga.service version.service xyz.openbmc_project.cpusensor.service
After=logrotate.service xyz.openbmc_project.Logging.IPMI.service fpga.service version.service xyz.openbmc_project.cpusensor.service

[Service]
Restart=always
RestartSec=5
StartLimitBurst=10
ExecStart=/usr/bin/fpgasensor

[Install]
WantedBy=multi-user.target
